# osso-af-utils/debian/postinst
#! /bin/sh
# postinst script for app-killer
#
# see: dh_installdeb(1)

set -e

if [ -e /targets/links/scratchbox.config ]; then
  . /targets/links/scratchbox.config
  USERHOME="$SBOX_TARGET_DIR/home/user"
else
  USERHOME="/home/user"
fi

# summary of how this script can be called:
#        * <postinst> `configure' <most-recently-configured-version>
#        * <old-postinst> `abort-upgrade' <new version>
#        * <conflictor's-postinst> `abort-remove' `in-favour' <package>
#          <new-version>
#        * <deconfigured's-postinst> `abort-deconfigure' `in-favour'
#          <failed-install-package> <version> `removing'
#          <conflicting-package> <version>
# for details, see http://www.debian.org/doc/debian-policy/ or
# the debian-policy package
#

case "$1" in
    configure)
      update-rc.d fb-progress.sh start 12 2 . stop 20 5 . > /dev/null
      USER_ICON_DIR=/usr/share/user-icons
      ln -sf $USER_ICON_DIR $USERHOME/.icons
      ln -sf /usr/share/icons/xcursor-transparent \
        $USER_ICON_DIR/xcursor-transparent
      ln -sf /usr/lib/libXcursor.so.1 /usr/lib/libXcursor.so
    ;;

    abort-upgrade|abort-remove|abort-deconfigure)

    ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0


