# osso-af-utils/debian/postrm
#!/bin/sh
update-rc.d -f fb-progress.sh remove > /dev/null
