# osso-af-startup/debian/postrm
#!/bin/sh

update-rc.d -f af-services remove > /dev/null
update-rc.d -f af-startup remove > /dev/null
update-rc.d -f af-base-apps remove > /dev/null
update-rc.d -f clean-tmp-dirs remove > /dev/null
